# files/flow.conf
#
# WARNING: http://startflowing.net/ server generated by ansible
#

description     "start and stop flow server"

start on runlevel [2345]
stop on runlevel [!2345]

# Give up if restart occurs 10 times in 30 seconds.
respawn
respawn limit 10 30

setgid jetty
setuid flow

chdir /opt/net/startflowing/flow/

env JAVA_HOME=/usr/lib/jvm/java-7-openjdk-amd64
env JAVA_OPTIONS="-Xmx1024m -XX:MaxPermSize=768m -XX:PermSize=256m -DFLOW_HOME=/var/lib/flow"
env JAVA_APP=jetty-runner-8.1.9.jar
env APP_OPTIONS="--lib /opt/jdbc/postgresql/9.3/ --port 8082 --path /flowui flowui-0.3.war"
env LOGFILE=/var/log/flow/jetty.log

script
  exec $JAVA_HOME/bin/java $JAVA_OPTIONS -jar $JAVA_APP $APP_OPTIONS > $LOGFILE 2>&1
end script

